# templates/init-api.j2
#!/bin/sh

# This file was generated by Ansible for {{ ansible_fqdn }}
# Do NOT modify this file by hand!

# uwsgi - Startup script for uwsgi
#
# chkconfig: - 85 15
# processname: gunicorn

# Source function library.
. /etc/rc.d/init.d/functions

# Also look at sysconfig; this is where environmental variables should be set
# on RHEL systems.
[ -f "/etc/sysconfig/graphite" ] && . /etc/sysconfig/graphite

prog="{{graphite_api_home}}/env/bin/gunicorn graphite_api.app:app"
user={{graphite_api_user}}
pidfile="/var/run/graphite.pid"
logfile="/var/log/graphite.log"
lockfile="/var/lock/subsys/graphite"
opts="\
  -b {{graphite_api_gunicorn_host}}:{{graphite_api_gunicorn_port}} \
  -w 2 \
  -u $user \
  -p $pidfile \
  -D \
"
export GRAPHITE_API_CONFIG={{ graphite_api_home }}/config.yml
RETVAL=0

start() {
    echo -n $"Starting graphite: "

    daemon --pidfile=${pidfile} ${prog} ${opts}
    RETVAL=$?
    echo
    [ $RETVAL = 0 ] && touch ${lockfile}
    return $RETVAL
}

stop() {
    echo -n $"Stopping graphite: "
    killproc -p ${pidfile} graphite
    RETVAL=$?
    echo
    [ $RETVAL = 0 ] && rm -f ${lockfile} ${pidfile}
}

reload() {
    echo -n $"Reloading graphite: "
    killproc -p ${pidfile} graphite -HUP
    RETVAL=$?
    echo
}

rh_status() {
    status -p ${pidfile} graphite
}

# See how we were called.
case "$1" in
    start)
        rh_status >/dev/null 2>&1 && exit 0
        start
        ;;
    stop)
        stop
        ;;
    status)
        rh_status
        RETVAL=$?
        ;;
    restart)
        stop
        start
        ;;
    force-reload|reload)
        reload
        ;;
    *)
        echo $"Usage: $0 {start|stop|restart|force-reload|reload|status}"
        RETVAL=2
esac

exit $RETVAL
